# easyeffects-auto-preset.service
[Unit]
Description=EasyEffects Auto Preset Switcher
After=pipewire.service pulseaudio.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=/home/enrico/easyeffects-auto-preset.sh start
Restart=always
RestartSec=3
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_RUNTIME_PATH=/run/user/1000/pulse

[Install]
WantedBy=default.target
